# etc/supervisor/conf.d/01-confd.conf
[program:confd]
command=/usr/local/sbin/confd -interval 10 -node 'http://127.0.0.1:4001' -config-file '/etc/confd/conf.d/haproxy.toml'
numprocs=1
autostart=true
autorestart=true
user=root
